# scripts/systemd/alpaca_hft_paper_start.service
[Unit]
Description=Start Alpaca HFT Paper Trader

[Service]
Type=oneshot
ExecStart=/bin/systemctl start alpaca_hft_paper.service
