# scripts/systemd/alpaca_hft_paper_stop.service
[Unit]
Description=Stop Alpaca HFT Paper Trader

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop alpaca_hft_paper.service
