# scripts/systemd/alpaca_hft_paper.service
[Unit]
Description=Alpaca HFT Paper Trader
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/Alphs
EnvironmentFile=/home/Alphs/.env
ExecStart=/usr/bin/python3 -m src.main run --strategies pairs,mm,leadlag
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
