# scripts/systemd/alpaca_hft_paper_start.timer
[Unit]
Description=Start Alpaca HFT Paper Trader at market open

[Timer]
OnCalendar=Mon..Fri 09:25:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
